# crontab
00 17 * * * python /pycode/main.py >> /logfile
